# access_daemon.service
[Unit]
Description=VPN Access Control Daemon
After=docker.service network.target
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/vpn-billing
ExecStart=/opt/vpn-billing/venv/bin/python3 /opt/vpn-billing/access_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
